# supervisord.conf
; ══════════════════════════════════════════════════════════════
;  OneIQ — Supervisord Configuration
;
;  BUG FIX: Removed "depends_on=php-fpm" — that is a Docker Compose
;  directive, NOT a valid supervisord key. It caused supervisord to
;  reject the [program:nginx] section entirely, leaving no HTTP
;  listener on $PORT → Railway health check → "service unavailable".
;
;  Ordering is handled via priority alone (php-fpm=1, nginx=2).
;  supervisord starts programs in ascending priority order.
; ══════════════════════════════════════════════════════════════

[supervisord]
nodaemon=true
; Use /tmp instead of /dev/stdout for the main supervisord log.
; Some supervisord builds on Alpine can't open /dev/stdout as a
; logfile due to permission issues when supervisord forks.
; Individual program logs still go to /dev/stdout / /dev/stderr.
logfile=/tmp/supervisord.log
logfile_maxbytes=10MB
logfile_backups=1
loglevel=info
pidfile=/tmp/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ── PHP-FPM (priority 1 — starts first) ───────────────────────
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F -R
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

; ── Nginx (priority 2 — starts after php-fpm) ─────────────────
; NOTE: Do NOT use "depends_on" here — that is Docker Compose syntax.
; Priority ordering ensures php-fpm starts first.
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2
